# nubis/bin/init-script-debian
#!/bin/bash

### BEGIN INIT INFO
# Provides:          confluence
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the confluence wiki system
# Description:       starts the confluence wiki system
### END INIT INFO

test -f /lib/lsb/init-functions && . /lib/lsb/init-functions

progname=`basename $0`

CONFLUENCE_START="/opt/atlassian/confluence/bin/startup.sh"
CONFLUENCE_STOP="/opt/atlassian/confluence/bin/shutdown.sh"
CONFLUENCE_USER="confluence"
RETVAL=0

# We assume that this will run on ubuntu
PREFIX='/opt/atlassian/confluence'
pidfile=/var/run/${CONFLUENCE_USER}/${progname}.pid

JAVA_HOME='/usr/lib/jvm/java-8-oracle'
CATALINA_HOME=${prefix}
CATALINA_PID=${pidfile}

export JAVA_HOME
export CATALINA_HOME

start_daemon() {
    export JAVA_HOME
    export CATALINA_PID
    export CATALINA_HOME

    # We assume /lib/lsb/init-functions is available
    log_daemon_msg "Starting ${progname}" ${progname}
    su - ${CONFLUENCE_USER} -c ${CONFLUENCE_START}
    RETVAL=$?
    if [[ $RETVAL == 0 ]]; then
        log_end_msg 0
    else
        log_end_msg 1
    fi
}

stop_daemon() {
    export JAVA_HOME
    export CATALINA_PID
    export CATALINA_HOME

    log_daemon_msg "Stopping ${progname}" ${progname}
    su - ${CONFLUENCE_USER} -c ${CONFLUENCE_STOP}
    RETVAL=$?
    if [[ $RETVAL == 0 ]]; then
        log_end_msg 0
    else
        log_end_msg 1
    fi
}

case "$1" in
    start)
        start_daemon
        ;;
    stop)
        stop_daemon
        ;;
    restart)
        stop_daemon
        sleep 5
        start_daemon
        ;;
    *)
        echo "Usage: confluence {start|stop|restart}"
        exit 1
        ;;
esac
exit $RETVAL
